# update-motd.d/10-uname
#!/bin/bash
c0='\e[1;31m' # green
c1='\e[1;32m' # red
c2='\e[1;34m' # blue

echo -e "
${c0}    __  __   
${c0}   (_\\)(/_) ${c2}                  ___   
${c1}   (_(__)_)  ${c2}   _______ _ __ / _ \ 
${c1}  (_(_)(_)_) ${c2}  |_  / _ \ '__| | | |
${c1}   (_(__)_)  ${c2}   / /  __/ |  | |_| |
${c1}     (__)    ${c2}  /___\___|_|   \___/ 
${c1}"
