# update-motd.d/60-tail
#! /bin/bash
if [ -f /etc/motd.tail ]; then
    cat /etc/motd.tail
fi
